# generic/all/custom-conf-sysv/usr/share/custom-conf-sysv/config/etc/rc.local
#!/bin/sh -e
## CUSTOM-CONF: File automatically created/updated [custom-conf-sysv]

# Run parts
run-parts --exit-on-error --verbose /etc/rc.local.d

# Done
exit 0

