# generic/all/custom-util-zfs/debian/prerm
#!/bin/bash
[ "${1}" != 'remove' ] && exit 0

# Stop the service
invoke-rc.d zfs-mount-late stop

# Exit
exit 0

